# fw-rules.cron
MAILTO=""
0 * * * * root /usr/share/firewall/turris
